# sentry-nightwatch/systemd/nightwatch-daily.timer
[Unit]
Description=Run NightWatch daily at 02:00 Asia/Tehran

[Timer]
OnCalendar=*-*-* 02:00:00 Asia/Tehran
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
